# contrib/init/bitcoind.conf
description "Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env GCOIND_BIN="/usr/bin/gcoind"
env GCOIND_USER="bitcoin"
env GCOIND_GROUP="bitcoin"
env GCOIND_PIDDIR="/var/run/gcoind"
# upstart can't handle variables constructed with other variables
env GCOIND_PIDFILE="/var/run/gcoind/gcoind.pid"
env GCOIND_CONFIGFILE="/etc/bitcoin/bitcoin.conf"
env GCOIND_DATADIR="/var/lib/gcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # gcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$GCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run gcoind."
        echo "The setting must appear in $GCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$GCOIND_PIDDIR"
    chmod 0755 "$GCOIND_PIDDIR"
    chown $GCOIND_USER:$GCOIND_GROUP "$GCOIND_PIDDIR"
    chown $GCOIND_USER:$GCOIND_GROUP "$GCOIND_CONFIGFILE"
    chmod 0660 "$GCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$GCOIND_PIDFILE" \
    --chuid $GCOIND_USER:$GCOIND_GROUP \
    --exec "$GCOIND_BIN" \
    -- \
    -pid="$GCOIND_PIDFILE" \
    -conf="$GCOIND_CONFIGFILE" \
    -datadir="$GCOIND_DATADIR" \
    -disablewallet \
    -daemon

